# alarmd/alarmd.service
[Unit]
Description=Alarm daemon, a software alarm clock

[Service]
ExecStart=$(PREFIX)/bin/alarmd

